# deploy/responsegrid.service
# ResponseGrid — reconcile the full prod stack on every boot.
#
# Why this exists: on start, the Docker daemon only re-applies each container's
# own `restart:` policy — it does NOT re-run `docker compose up`. The `migrate`
# service is one-shot (restart: no) and `api` depends on it
# (service_completed_successfully), so after an EC2 stop/start the daemon brings
# back postgres/redis/caddy/datadog but NOT migrate, and `api` therefore stays
# down → Caddy returns 502 until someone runs `docker compose up -d` by hand.
#
# This oneshot unit runs exactly that command once, after Docker is up, on every
# boot. It reconciles the compose desired state: migrate re-runs (idempotent),
# then api starts. RemainAfterExit keeps the unit "active (exited)" afterwards.
#
# Source of truth for the unit. Installed by:
#   - deploy/user-data.sh           (fresh instances — written inline, repo not cloned yet)
#   - deploy/install-reconcile.sh   (existing instances + every CI deploy — copies THIS file)
# Keep the inline copy in deploy/user-data.sh in sync with this file.
#
# Manual install (from the repo root, e.g. /opt/responsegrid):
#   sudo bash deploy/install-reconcile.sh

[Unit]
Description=ResponseGrid stack (docker compose up -d) — reconcile on every boot
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
# Skip cleanly (not fail) on a brand-new box before the repo/.env exist.
ConditionPathExists=/opt/responsegrid/deploy/docker-compose.prod.yml
ConditionPathExists=/opt/responsegrid/deploy/.env

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/responsegrid
ExecStart=/usr/bin/docker compose -f deploy/docker-compose.prod.yml up -d
# Cold boot needs to wait for postgres health + the one-shot migrate to finish
# (compose honours depends_on conditions even with -d), so allow generous time.
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
